# vendor/github.com/mesosphere/dcos-edge-lb/apiserver/containerfiles/apiserver/service/apiserver/run
#!/bin/bash
set -euo pipefail

log() {
    msg="$1"

    logline="[$(pwd) $0] $msg\n"

    printf "$logline" >&1
    printf "$logline" >&2
}

log "starting"

# Local
: ${APIS_SCHEME:="http"}
: ${PORT0:=52978}
: ${PORT1:=52979}
: ${APIS_HOST:="0.0.0.0"}
: ${APIS_PORT:=$PORT0}
: ${APIS_TLS_HOST:="0.0.0.0"}
: ${APIS_TLS_PORT:=$PORT1}
: ${APIS_TLS_CERTIFICATE:=""}
: ${APIS_TLS_PRIVATE_KEY:=""}
: ${APIS_TLS_CA_CERTIFICATE:=""}
: ${DCOS_SERVICE_ACCOUNT_CREDENTIAL:=""}
: ${DCOS_SECRET_NAME:=""}
: ${DCOS_PRINCIPAL:=""}
: ${APIS_ZK:="zk-1.zk:2181,zk-2.zk:2181,zk-3.zk:2181,zk-4.zk:2181,zk-5.zk:2181"}
: ${APIS_ZK_TIMEOUT:="60s"}
: ${APIS_ZK_PATH:="/apiserver"}
: ${ELB_CFGCACHE_FILENAME:="cfgcache.json"}
: ${ELB_CFGCACHE_FILE:="$LBWORKDIR/$ELB_CFGCACHE_FILENAME"}
: ${DCOS_TEMPLATE_SV:="$LBWORKDIR/apiserver/service/dcos-template"}

# Exported
export ELB_CFGCACHE_FILENAME
export ELB_CFGCACHE_FILE
export DCOS_TEMPLATE_SV
export DCOS_SERVICE_ACCOUNT_CREDENTIAL
export DCOS_SECRET_NAME
export DCOS_PRINCIPAL
export APIS_ZK
export APIS_ZK_TIMEOUT
export APIS_ZK_PATH
export HOST=$APIS_HOST
export PORT=$APIS_PORT
export TLS_HOST=$APIS_TLS_HOST
export TLS_PORT=$APIS_TLS_PORT
export TLS_CERTIFICATE=$APIS_TLS_CERTIFICATE
export TLS_PRIVATE_KEY=$APIS_TLS_PRIVATE_KEY
export TLS_CA_CERTIFICATE=$APIS_TLS_CA_CERTIFICATE

exec "$LBWORKDIR/apiserver/bin/apiserver" --scheme $APIS_SCHEME --verbose
